# supervisord.conf
[supervisord]
nodaemon=true

[program:startup]
command=/bin/bash /startup.sh
user=root
autostart=true
autorestart=true
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
